# etc/rc.local
#!/bin/sh -e

ptpd2 --slaveonly --interface eth0

exit 0
